# supervisord.conf
; supervisor config file for Voice AI Platform
; RunPod deployment with multiple services

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ==================== GPT-SoVITS Workers ====================
; Using official api.py with different ports
; Each worker needs its own config or uses default parameters

[program:gpt-sovits-1]
command=python api.py -a 127.0.0.1 -p 9881 -dr "/workspace/models/custom_voices/default/reference.wav" -dt "参考音频文本" -dl "zh"
directory=/workspace/GPT-SoVITS
autostart=false
autorestart=true
priority=10
startsecs=30
stdout_logfile=/var/log/supervisor/gpt-sovits-1.log
stderr_logfile=/var/log/supervisor/gpt-sovits-1-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

[program:gpt-sovits-2]
command=python api.py -a 127.0.0.1 -p 9882 -dr "/workspace/models/custom_voices/default/reference.wav" -dt "参考音频文本" -dl "zh"
directory=/workspace/GPT-SoVITS
autostart=false
autorestart=true
priority=10
startsecs=30
stdout_logfile=/var/log/supervisor/gpt-sovits-2.log
stderr_logfile=/var/log/supervisor/gpt-sovits-2-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

[program:gpt-sovits-3]
command=python api.py -a 127.0.0.1 -p 9883 -dr "/workspace/models/custom_voices/default/reference.wav" -dt "参考音频文本" -dl "zh"
directory=/workspace/GPT-SoVITS
autostart=false
autorestart=true
priority=10
startsecs=30
stdout_logfile=/var/log/supervisor/gpt-sovits-3.log
stderr_logfile=/var/log/supervisor/gpt-sovits-3-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

[program:gpt-sovits-4]
command=python api.py -a 127.0.0.1 -p 9884 -dr "/workspace/models/custom_voices/default/reference.wav" -dt "参考音频文本" -dl "zh"
directory=/workspace/GPT-SoVITS
autostart=false
autorestart=true
priority=10
startsecs=30
stdout_logfile=/var/log/supervisor/gpt-sovits-4.log
stderr_logfile=/var/log/supervisor/gpt-sovits-4-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

; ==================== RVC Workers ====================
; Using FastAPI wrapper (SocAIty fork or custom api.py)

[program:rvc-1]
command=python api.py --port 7866 --device cuda
directory=/workspace/RVC
autostart=false
autorestart=true
priority=20
startsecs=20
stdout_logfile=/var/log/supervisor/rvc-1.log
stderr_logfile=/var/log/supervisor/rvc-1-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

[program:rvc-2]
command=python api.py --port 7867 --device cuda
directory=/workspace/RVC
autostart=false
autorestart=true
priority=20
startsecs=20
stdout_logfile=/var/log/supervisor/rvc-2.log
stderr_logfile=/var/log/supervisor/rvc-2-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

; ==================== Gateway & WebUI ====================
; Start after workers are ready

[program:model-webui]
command=python /workspace/webui/app.py
directory=/workspace
autostart=true
autorestart=true
priority=50
startsecs=10
stdout_logfile=/var/log/supervisor/webui.log
stderr_logfile=/var/log/supervisor/webui-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

[program:api-gateway]
command=python /workspace/gateway/main.py
directory=/workspace
autostart=true
autorestart=true
priority=100
startsecs=5
stdout_logfile=/var/log/supervisor/api-gateway.log
stderr_logfile=/var/log/supervisor/api-gateway-error.log
environment=PYTHONUNBUFFERED="1",CUDA_VISIBLE_DEVICES="0"

; ==================== Group Definitions ====================
[group:gpt-sovits]
programs=gpt-sovits-1,gpt-sovits-2,gpt-sovits-3,gpt-sovits-4

[group:rvc]
programs=rvc-1,rvc-2
